# s6/nginx/finish
#!/bin/sh

## Add your cleanup tasks above this line
echo "SERVICE ENDED: nginx-service"